# 2013_08_01/latest_build/package/rmqctl/src/mconfread
#!/bin/sh /etc/rc.common
. /etc/functions.sh

#START=40
#STOP=90
config_clear
config_load mcast
ENABLE_MCAST=0
HAS_TUNNEL=0

load_src() {
    local channel
    local tunnel
    local localip
    local remoteip
    local srcip
    local grpip  

    config_get channel $1 channel
    config_get tunnel $1 tunnel
    config_get localip $1 localip
    config_get remoteip $1 remoteip
    config_get srcip $1 srcip
    config_get grpip $1 grpip
 
    [ -z "$channel" ] && {
        [ "$tunnel" = "yes" ] && {                
            [ -n "$localip" -a -n "$remoteip" ] && {
                HAS_TUNNEL=1
                echo "tunnel $localip $remoteip" >> /var/run/mrouted.conf
            }
        }
    }   

    if [ -n "$channel" ] ; then {
        echo iptables -I zone_wan -p 6 -j ACCEPT >> /var/run/mcast.fw
        echo iptables -I zone_lan -p 6 -j ACCEPT >> /var/run/mcast.fw
    } else { 
        [ "$tunnel" = "yes" ] && {  
            echo iptables -I zone_wan -p 4 -j ACCEPT >> /var/run/mcast.fw
            echo iptables -I zone_lan -p 4 -j ACCEPT >> /var/run/mcast.fw            
        }
    }
    fi      

    [ -n "$srcip" -a -n "$grpip" ] && {
        echo iptables -I zone_lan_forward -p udp -s $srcip -d $grpip -j ACCEPT >> /var/run/mcast.fw
        echo iptables -I zone_wan_forward -p udp -s $srcip -d $grpip -j ACCEPT >> /var/run/mcast.fw 
        echo iptables -I zone_lan -p udp -s $srcip -d $grpip -j ACCEPT >> /var/run/mcast.fw    
        echo iptables -I zone_wan -p udp -s $srcip -d $grpip -j ACCEPT >> /var/run/mcast.fw
    }             
}

load_rate() {
    local iface
    local rate

    config_get iface $1 iface
    config_get rate $1 rate
    [ -n "$iface" -a -n "$rate" ] && {
        testmode $iface mcast set $rate
    }
}

check_able() {
    local enable    
    config_get enable $1 enable
    ENABLE_MCAST=$enable 
}

start() {
    echo "Setting up multicast rate"
    config_foreach load_rate mrate

    sysmode=$(uci get system.@system[0].systemmode)
    if [ "$sysmode" = "router" -o "$sysmode" = "gateway" ] ; then
        config_foreach check_able mcontrol
        [ "$ENABLE_MCAST" = "1" ] && {
            echo "Loading multicast src"        
            mkdir -p /var/run
            rm -f /var/run/mcast.fw
            rm -f /var/run/mrouted.conf
            echo "#Don't modify, this is generated by mconfread" > /var/run/mrouted.conf
            echo iptables -I zone_lan_forward -p igmp -j ACCEPT >> /var/run/mcast.fw
            echo iptables -I zone_wan_forward -p igmp -j ACCEPT >> /var/run/mcast.fw
            echo iptables -I zone_lan -p igmp -j ACCEPT >> /var/run/mcast.fw    
            echo iptables -I zone_wan -p igmp -j ACCEPT >> /var/run/mcast.fw
            echo iptables -I zone_lan -p udp -d 224.0.0.9 -j ACCEPT  >> /var/run/mcast.fw
            echo iptables -I zone_wan -p udp -d 224.0.0.9 -j ACCEPT  >> /var/run/mcast.fw
            echo iptables -I zone_lan_forward -p 4 -j ACCEPT >> /var/run/mcast.fw
            echo iptables -I zone_wan_forward -p 4 -j ACCEPT >> /var/run/mcast.fw

            config_foreach load_src mcast
            ln -s /var/run/mcast.fw /etc/mcast.fw
            mrouted -c /var/run/mrouted.conf
            [ "$HAS_TUNNEL" = "1" ] && {
                while [ -z "$DVMRPS" ] ; do
                    sleep 1
                    DVMRPS=$(cat /var/run/tunlfl)                
                done
            }  
            /etc/init.d/firewall restart
            return  
        }
    fi
    stop
}


stop() {
    rm -rf /etc/mcast.fw 
    killall -9 mrouted 1>>/dev/null 2>>/dev/null
}

restart() {
    stop
    start
}


