# 2013_08_01/latest_build/package/base-files/files/etc/hotplug.d/CVS/Entries
D/button////
D/ieee1394////
D/iface////
D/net////
D/usb////
