# 2013_08_01/latest_build/package/websockserver/files/websockserverd.init
#!/bin/sh /etc/rc.common

. /etc/version

START=99
STOP=99

start(){
        killall -9 websockserver
        killall -9 websockscript
        /bin/websockscript &
}
stop(){
        killall -9 websockserver
        killall -9 websockscript
}

