# 2013_08_01/latest_build/package/msmscomd/files/salertd.init
#!/bin/sh /etc/rc.common

START=99 
STOP=99

start(){
  killall -9 salert
  /bin/salert >/dev/null 2>&1 &
}

stop(){
  killall -9 salert
}


